# nginx/laravel-worker.conf
[program:laravel-worker]
command=php /Data/artisan queue:work --sleep=3 --tries=3 --daemon
autostart=true
autorestart=true
redirect_stderr=true